# scripts/minecraft-shutdown.service
[Unit]
Description=Minecraft Server Graceful Shutdown
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/graceful-shutdown.sh
TimeoutStartSec=600
RemainAfterExit=yes

[Install]
WantedBy=shutdown.target reboot.target halt.target